# deploy/mangastore.service
[Unit]
Description=MangaStore Flask Application
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mangastore
Environment=FLASK_ENV=production
ExecStart=/home/pi/mangastore/venv/bin/python run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
